# borg-backup-pictures.timer
[Unit]
Description=Borg for Pictures

[Timer]
OnBootSec=15m
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target