# borg-backup-tuwien.timer
[Unit]
Description=Borg Backup: TUWien

[Timer]
OnCalendar=*-*-* 00..02,12,20..23:00:00

[Install]
WantedBy=timers.target