# borg-backup-libvirt.timer
[Unit]
Description=Borg for Pictures

[Timer]
OnBootSec=1h
OnActiveSec=4h

[Install]
WantedBy=timers.target