# borg-backup-libvirt.timer
[Unit]
Description=Borg Backup: libvirt

[Timer]
OnCalendar=*-*-* 00,12,16,22:00:00

[Install]
WantedBy=timers.target