# borg-backup-documents.timer
[Unit]
Description=Automated Borg Backup: Documents

[Timer]
WakeSystem=false

OnBootSec=15min
OnActiveSec=15min

[Install]
WantedBy=timers.target