# borg-backup-dotfiles.timer
[Unit]
Description=Automated Borg Backup: dotfiles

[Timer]
OnBootSec=15m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target